# scripts/linux/desktop-commander.socket
[Unit]
Description=DesktopCommanderMCP Unix Socket
Documentation=https://github.com/wonderwhy-er/DesktopCommanderMCP
After=network.target

[Socket]
ListenStream=/run/desktop-commander/socket
SocketMode=0666
DirectoryMode=0755
Accept=no
MaxConnections=100
MaxConnectionsPerSource=10

[Install]
WantedBy=sockets.target
