# scripts/linux/desktop-commander.service
[Unit]
Description=DesktopCommanderMCP Unix Socket Handler
Documentation=https://github.com/wonderwhy-er/DesktopCommanderMCP
After=network.target desktop-commander.socket
Requires=desktop-commander.socket

[Service]
Type=simple
ExecStart=/usr/bin/node /persistent/home/raka/mcp-servers/DesktopCommanderMCP/dist/index.js --socket-mode
WorkingDirectory=/persistent/home/raka/mcp-servers/DesktopCommanderMCP
Restart=always
RestartSec=5
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true

# Allow writing to socket directory
ReadWritePaths=/run/desktop-commander

# Environment
Environment=NODE_ENV=production
Environment=DESKTOP_COMMANDER_MODE=unix_socket

[Install]
WantedBy=multi-user.target
